# templates/default/tomcat-initscript.sh.erb
#!/bin/bash

# !!!!!!! This file created by Chef !!!!!!!
# Any edits to this file will be overwritten

## $Id$

#
# tomcat	Starts/stops the tomcat service
#
# chkconfig: 35 80 20
# description: Tomcat application server

### BEGIN INIT INFO
# Provides: tomcat
# Required-Start: $local_fs $network $syslog
# Required-Stop: $local_fs $syslog
# Should-Start: $syslog
# Should-Stop: $network $syslog
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: Tomcat application server
# Description:       Tomcat application server
### END INIT INFO
export JAVA_HOME=<%= @java_home %>
/bin/su -s /bin/bash <%= @tomcat_user %> -- <%= @tomcat_home %>/bin/tomcat "$@"
